# deploy/systemd/echochat@.service
# Echo-Chat systemd template for multiple one-worker instances
#
# Preferred path:
#   python main.py --write-deployment-kit
#
# Install example:
#   sudo cp deploy/systemd/echochat@.service /etc/systemd/system/echochat@.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now echochat@5000 echochat@5001
#
# The instance name is the backend port. Keep every process at one Gunicorn
# worker. For 10 instances use ports such as 5000-5009 behind sticky reverse-
# proxy routing and configure ECHOCHAT_SOCKETIO_MESSAGE_QUEUE=redis://127.0.0.1:6379/1.

[Unit]
Description=Echo-Chat one-worker instance on port %i
After=network-online.target redis.service
Wants=network-online.target redis.service

[Service]
Type=simple
WorkingDirectory=/opt/echochat/Echo-Chat-main
EnvironmentFile=/etc/echochat/echochat.env
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/opt/echochat/Echo-Chat-main/.venv/bin/python tools/config_doctor.py --config /opt/echochat/Echo-Chat-main/server_config.json
ExecStartPre=/opt/echochat/Echo-Chat-main/.venv/bin/python main.py --config /opt/echochat/Echo-Chat-main/server_config.json --redis-socketio-check
# Bind overrides live in ExecStart so EnvironmentFile cannot force every instance back to port 5000.
ExecStart=/usr/bin/env ECHOCHAT_CONFIG=/opt/echochat/Echo-Chat-main/server_config.json ECHOCHAT_BIND=127.0.0.1:%i ECHOCHAT_PRODUCTION_BIND=127.0.0.1:%i ECHOCHAT_PRODUCTION_WORKERS=1 ECHOCHAT_WORKERS=1 /opt/echochat/Echo-Chat-main/.venv/bin/python main.py --production
Restart=on-failure
RestartSec=3
TimeoutStopSec=25

User=echochat
Group=echochat

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=full
RestrictSUIDSGID=true
LockPersonality=true

ReadWritePaths=/opt/echochat/Echo-Chat-main/logs
ReadWritePaths=/opt/echochat/Echo-Chat-main/uploads
ReadWritePaths=/opt/echochat/Echo-Chat-main/private_uploads
ReadWritePaths=/opt/echochat/Echo-Chat-main/instance
ReadWritePaths=/opt/echochat/Echo-Chat-main/static/uploads
ReadWritePaths=/opt/echochat/Echo-Chat-main/server_config.json

[Install]
WantedBy=multi-user.target
