# deploy/systemd/echochat-janitor.service
# Echo-Chat janitor service
#
# Run exactly one cleanup process when the chat server is served by Gunicorn or
# multiple one-worker instances. Do not start one janitor per web instance.
#
# Preferred path:
#   python main.py --write-deployment-kit
#
# Install:
#   sudo cp deploy/systemd/echochat-janitor.service /etc/systemd/system/echochat-janitor.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now echochat-janitor

[Unit]
Description=Chat server janitor for Echo-Chat deployments
After=network-online.target redis.service
Wants=network-online.target redis.service

[Service]
Type=simple
WorkingDirectory=/opt/echochat/Echo-Chat-main
EnvironmentFile=/etc/echochat/echochat.env
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/opt/echochat/Echo-Chat-main/.venv/bin/python tools/config_doctor.py --config /opt/echochat/Echo-Chat-main/server_config.json
ExecStart=/opt/echochat/Echo-Chat-main/.venv/bin/python janitor_runner.py --config /opt/echochat/Echo-Chat-main/server_config.json
Restart=on-failure
RestartSec=3
TimeoutStopSec=25

User=echochat
Group=echochat

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=full
RestrictSUIDSGID=true
LockPersonality=true

ReadWritePaths=/opt/echochat/Echo-Chat-main/logs
ReadWritePaths=/opt/echochat/Echo-Chat-main/uploads
ReadWritePaths=/opt/echochat/Echo-Chat-main/private_uploads
ReadWritePaths=/opt/echochat/Echo-Chat-main/instance
ReadWritePaths=/opt/echochat/Echo-Chat-main/static/uploads

[Install]
WantedBy=multi-user.target
